# server/agora.conf
start on startup
stop on shutdown
respawn
respawn limit 20 5

script
  exec /usr/local/bin/node /home/ec2-user/server/server.js
end script

post-start script
  echo "Node server started"
end script
